# deploy/systemd/bom-sync.timer
[Unit]
Description=PLM Data Sync Timer (Daily at 2 AM)
Requires=bom-sync.service

[Timer]
# 每天凌晨2点执行
OnCalendar=*-*-* 02:00:00
# 如果错过了执行时间，立即执行
Persistent=true
# 首次启动后延迟5分钟执行（避免系统启动时立即执行）
OnBootSec=5min

[Install]
WantedBy=timers.target
